# ai-factory/ops/com.ai-factory.poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-factory.poller</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/senioraiconsultant/.local/bin/npx</string>
    <string>tsx</string>
    <string>src/sources/poll-linear.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/senioraiconsultant/Development/Edu/ai-sdlc/ai-factory</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/senioraiconsultant/.local/bin:/Users/senioraiconsultant/.kimi-code/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>15</integer>

  <key>StandardOutPath</key>
  <string>/Users/senioraiconsultant/.ai-factory/logs/poller.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/senioraiconsultant/.ai-factory/logs/poller.err.log</string>
</dict>
</plist>
